# health-reminder.service
[Unit]
Description=Health Reminder Daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
ExecStart=/usr/local/bin/health-reminder start
Restart=on-failure
RestartSec=5
Environment=QT_QPA_PLATFORM=xcb

[Install]
WantedBy=default.target
